# going-vimsane.service
[Unit]
Description=Going Vimsane - Kanata layer visual overlay
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/luar/Workspace/projects/desktop/going-vimsane/going-vimsane
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
